# launch/pmx_system.service
[Unit]
After=network-online.target
Description=ROS2 PIPM services

[Service]
User=putm
Group=putm
ExecStart=/bin/bash -c 'source /opt/ros/humble/setup.bash; source /home/putm/putm_ws/install/setup.bash; ros2 launch /home/putm/putm_ws/launch/putm_ws.launch.py;'
RemainAfterExit=no
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=default.target